# server/systemd/cron_units/heartbeat.service
[Unit]
Description=Antigravity Heartbeat — system health check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=antigravity-ai
WorkingDirectory=/home/antigravity-ai/Antigravity_Orchestra
ExecStart=/home/antigravity-ai/Antigravity_Orchestra/run_cron.sh orchestra_heartbeat.py
TimeoutStartSec=300
StandardOutput=append:/home/antigravity-ai/Antigravity_Orchestra/logs/heartbeat.log
StandardError=append:/home/antigravity-ai/Antigravity_Orchestra/logs/heartbeat.log
