# server/systemd/cron_units/research-daily.service
[Unit]
Description=Antigravity Scholar — daily research briefs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=antigravity-ai
WorkingDirectory=/home/antigravity-ai/Antigravity_Orchestra
ExecStart=/home/antigravity-ai/Antigravity_Orchestra/run_cron.sh research_engine.py --daily-cron
TimeoutStartSec=600
StandardOutput=append:/home/antigravity-ai/Antigravity_Orchestra/logs/research-daily.log
StandardError=append:/home/antigravity-ai/Antigravity_Orchestra/logs/research-daily.log
