# server/systemd/cron_units/social-listener.service
[Unit]
Description=Antigravity Social Listener — polls chatroom for trigger tags
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=antigravity-ai
WorkingDirectory=/home/antigravity-ai/Antigravity_Orchestra
ExecStart=/home/antigravity-ai/Antigravity_Orchestra/run_cron.sh social_listener_serverside.py
TimeoutStartSec=300
StandardOutput=append:/home/antigravity-ai/Antigravity_Orchestra/logs/social-listener.log
StandardError=append:/home/antigravity-ai/Antigravity_Orchestra/logs/social-listener.log
