# server/systemd/cron_units/social-listener.timer
[Unit]
Description=Social Listener — every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
